# webapp/salt/roots/gowebapp.conf
[program:gowebapp]
command=export PATH="$PATH:/usr/local/go/bin" ; go run gowebapp.go
autostart=true
autorestart=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log
